# kmu-test.service
[Unit]
Description=KMU-test
[Service]
User=nobody
ExecStart=/usr/bin/go_build_main